# etc/cron.d/nsm-watchdog
# /etc/cron.d/nsm-watchdog
#
# crontab entry to restart processes if they fail

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

4,9,14,19,24,29,34,39,44,49,54,59 *   * * * root  /usr/sbin/so-nsm-watchdog >> /var/log/nsm/watchdog.log 2>&1
