# etc/cron.d/server-backup-config
# /etc/cron.d/server-backup-config
#
# crontab entry to backup server config; purge backups older than 10 days
DAYSTOKEEP=10
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
0 1 * * * root if [ -d /etc/nsm/backup ]; then :; else mkdir /etc/nsm/backup; fi; /usr/sbin/nsm_server_backup-config --force-yes --server-name=securityonion --backup-file=/etc/nsm/backup/securityonion-server-backup-`date +\%Y-\%m-\%d`.tar.gz >> /var/log/nsm/server-backup-config.log 2>&1; /usr/bin/find  /etc/nsm/backup/*server-backup*.tar.gz -mtime +$DAYSTOKEEP -exec rm -f {} \;

