# etc/cron.d/so-server-backup-config
# /etc/cron.d/so-server-backup
#
# crontab entry to backup server config
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
0 1 * * * root /usr/sbin/so-server-backup
