# etc/cron.d/so-sensor-backup-config
# /etc/cron.d/so-sensor-backup-config
#
# crontab entry to backup sensor config
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
0 1 * * * root /usr/sbin/so-sensor-backup-config
